# system/multiband-webapp.service
[Service]
	After=multiband-modulemanager.service
	ExecStart=/home/andrew/MultiBand/WebApp/webserver.sh
	Restart=always
	StandardOutput=syslog
	StandardError=syslog
	SyslogIdentifier=mb-webapp
	User=andrew
[Install]